# install/ubuntu/redwood-router.conf
exec /var/www/redwood/go/bin/expecon-router
start on runlevel [23]
